# cron.d/update-oS-helper.sh
#!/bin/bash

pushd /home/dimstar/src/openSUSE-helpers
git pull

